# server/supervisor/qux.conf
[program:qux]
command = node app.js
directory = /www/qux-tool/server
stopsignal = QUIT
user = root
autostart = true
autorestart = true
